# qbittorrent/clean-data-Downloads.service
[Unit]
Description=clean-data-Downloads

[Service]
ExecStart=/usr/bin/clean-data-Downloads

[Install]
WantedBy=multi-user.target
