# qbittorrent/clean-data-Downloads.service
[Unit]
Description=clean-data-Downloads

[Service]
ExecStart=bash -xc 'rm -rv /mnt/data/Downloads/*'

[Install]
WantedBy=multi-user.target
