# qbittorrent/openvpn-reconnect.service
[Unit]
Description=openvpn-reconnect

[Service]
ExecStart=/usr/bin/openvpn-reconnect

[Install]
WantedBy=multi-user.target
